# twitter-analyzer.service
[Unit]
Description=Twitter Analyzer
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/twitter-analyzer
EnvironmentFile=/opt/twitter-analyzer/.env
ExecStart=/opt/twitter-analyzer/venv/bin/gunicorn -w 4 -k uvicorn.workers.UvicornWorker app:app --bind 0.0.0.0:8000
Restart=always

[Install]
WantedBy=multi-user.target
